# rootfs/etc/crontabs/www-data
* * * * * php /var/www/wordpress/wp-cron.php
