# rootfs/service/nginx/finish
#!/bin/sh
sv exit /service/*
kill -HUP $(pidof runsvdir)
