# rootfs/service/nginx/run
#!/usr/bin/dumb-init /bin/sh
sv -w5 check php-fpm
exec nginx
