# etc/supervisor.d/postfix.ini
[program:postfix]
command=start-postfix.sh
user=root
autostart=true
autorestart=false
priority=10
