# rootfs/etc/service/nginx/run
#!/bin/bash -e

# Wait until syslog is up
until [ -S /dev/log ]; do sleep 1; done

exec nginx
