# rootfs/service/php-fpm/run
#!/usr/bin/dumb-init /bin/sh
exec chpst -e /etc/environment php-fpm
